# DevOps/etc/profile.d/MYSQL_HOME.sh
export MYSQL_HOME=/opt/mysql
export PATH=$MYSQL_HOME/bin:$PATH