# DevOps/etc/profile.d/NODE_HOME.sh
export NODE_HOME=/opt/node-v14.15.4-linux-x64
export PATH=$NODE_HOME/bin:$PATH