# DevOps/etc/profile.d/SONAR_SCANNER_HOME.sh
export SONAR_SCANNER_HOME=/opt/sonar-scanner-4.3.0.2102-linux
export PATH=$SONAR_SCANNER_HOME/bin:$PATH